# scripts/sensors.service
[Unit]
Description=Fetch HA sensor data cache
After=network-online.target

[Service]
Type=oneshot
User=sethpthomas91
ExecStart=/var/www/homehub/scripts/fetch_sensors.sh
