# scripts/weather.service
[Unit]
Description=Fetch Open-Meteo weather cache
After=network-online.target

[Service]
Type=oneshot
User=sethpthomas91
ExecStart=/var/www/homehub/scripts/fetch_weather.sh
