# scripts/weather.timer
[Unit]
Description=Refresh weather cache every 15 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
